# macos/local.proxy-detect.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- proxy-config-utility — macOS LaunchAgent template (automatic detection).
     Runs proxy-detect at login, every 5 minutes, and on network change, so your
     shells' cache stays fresh with no per-shell work.

     Install:
       1. CUSTOMIZE the ProgramArguments path below to your absolute proxy-detect
          path (launchd does not expand ~; e.g. /Users/you/bin/proxy-detect), and
          the two log paths.
       2. cp macos/local.proxy-detect.plist ~/Library/LaunchAgents/
       3. launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/local.proxy-detect.plist
          launchctl kickstart -k gui/$(id -u)/local.proxy-detect
     Uninstall:
          launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/local.proxy-detect.plist

     StartInterval=300 is a safety net if a WatchPaths notification is missed.
     ThrottleInterval coalesces network-flap bursts (the mkdir lock makes overlap
     harmless anyway). -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>local.proxy-detect</string>
  <key>ProgramArguments</key><array><string>/ABSOLUTE/PATH/TO/bin/proxy-detect</string></array>
  <key>RunAtLoad</key><true/>
  <key>StartInterval</key><integer>300</integer>
  <key>WatchPaths</key><array>
    <string>/etc/resolv.conf</string>
    <string>/var/run/resolv.conf</string>
    <string>/Library/Preferences/SystemConfiguration</string>
  </array>
  <key>ThrottleInterval</key><integer>10</integer>
  <key>StandardOutPath</key><string>/ABSOLUTE/PATH/TO/.cache/proxy/launchd.log</string>
  <key>StandardErrorPath</key><string>/ABSOLUTE/PATH/TO/.cache/proxy/launchd.log</string>
  <key>ProcessType</key><string>Background</string>
</dict></plist>
